# config/launchd/com.axion.core.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Axion Core launchd template.
  NOTE: This is a TEMPLATE. The actual plist is generated by axion-launcher
  or install-mac.sh with correct absolute paths. If using this template
  directly, replace AXION_PROJECT_DIR and AXION_DATA_DIR with absolute paths
  (launchd does NOT expand ~ in plist values).
-->
<dict>
    <key>Label</key>
    <string>com.axion.core</string>

    <key>ProgramArguments</key>
    <array>
        <string>AXION_PROJECT_DIR/.venv/bin/python</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>src.main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>7777</string>
        <string>--log-level</string>
        <string>info</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>WorkingDirectory</key>
    <string>AXION_PROJECT_DIR</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>AXION_PROJECT_DIR/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>KLEITOS_DATA_DIR</key>
        <string>AXION_DATA_DIR</string>
        <key>KLEITOS_DB_PATH</key>
        <string>AXION_DATA_DIR/db/kleitos.db</string>
    </dict>

    <key>StandardOutPath</key>
    <string>AXION_DATA_DIR/logs/kleitos-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>AXION_DATA_DIR/logs/kleitos-stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
